# packages/wol-forwarder/data/etc/init.d/wol-forwarder
#!/bin/sh /etc/rc.common

START=99
EXTRA_COMMANDS="status"
EXTRA_HELP="        status          Service status"

config_get() {
  uci get "$1" 2>/dev/null || echo "$2"
}

start() {
  is_running && stop && sleep 1
  WOL_ADDR="$(config_get "wol-forwarder.main.addr" "0.0.0.0")"
  WOL_PORT="$(config_get "wol-forwarder.main.port" "1999")"
  WOL_BADDR="$(config_get "wol-forwarder.main.baddr" "255.255.255.255")"
  WOL_BPORT="$(config_get "wol-forwarder.main.bport" "9")"
  /usr/bin/wol-forwarder \
    -addr "${WOL_ADDR}" \
    -port "${WOL_PORT}" \
    -baddr "${WOL_BADDR}" \
    -bport "${WOL_BPORT}" \
    >/var/log/wol-forwarder.log &
}

stop() {
  killall /usr/bin/wol-forwarder >/dev/null 2>&1
}

restart() {
  start
}

reload() {
  restart
}

is_running() {
  pidof /usr/bin/wol-forwarder >/dev/null
}

status() {
  is_running && echo "running" || echo "inactive" && false
}
