# packages/lego/data/etc/init.d/lego
#!/bin/sh /etc/rc.common

START=21
EXTRA_COMMANDS="status"
EXTRA_HELP="        status          Service status"

. /etc/lego/env

remove_jobs() {
  sed -i '/lego_wrapper/d' /etc/crontabs/root >/dev/null 2>&1
}

start() {
  remove_jobs
  [ ! -f /etc/crontabs/root ] && mkdir -p /etc/crontabs && touch /etc/crontabs/root
  /usr/bin/lego_wrapper -a -f "${LEGO_PATH}/domains" >>/var/log/lego.log 2>&1 &
  echo "${LEGO_CRON_CONF} /usr/bin/lego_wrapper -f \"${LEGO_PATH}/domains\" >>/var/log/lego.log 2>&1" >>/etc/crontabs/root
  if [ "$(/etc/init.d/cron status)" = "running" ]; then
    /etc/init.d/cron reload
  fi
}

stop() {
  remove_jobs
  if [ "$(/etc/init.d/cron status)" = "running" ]; then
    /etc/init.d/cron reload
  fi
}

restart() {
  start
}

reload() {
  restart
}

status() {
  if [ "$(/etc/init.d/cron status)" != "running" ]; then
    echo "inactive" && false
  else
    grep -q lego_wrapper /etc/crontabs/root && echo "running" || echo "inactive" && false
  fi
}
