# packages/lego/data/etc/init.d/lego
#!/bin/sh /etc/rc.common

START=21
EXTRA_COMMANDS="status"
EXTRA_HELP="        status          Service status"

. /etc/lego/env

remove_jobs() {
  sed -i '/lego_wrapper/d' /etc/crontabs/root >/dev/null 2>&1
}

start() {
  remove_jobs
  while IFS= read -r domain; do
    domain_file_name="$(echo "${domain}" | sed "s/\*/_/g")"
    /usr/bin/lego_wrapper "${domain}" >>/var/log/lego_${domain_file_name}.log 2>&1 &
    echo "${LEGO_CRON_CONF} /usr/bin/lego_wrapper \"${domain}\" >>/var/log/lego_${domain_file_name}.log 2>&1" >>/etc/crontabs/root
  done </etc/lego/domains
  /etc/init.d/cron reload
}

stop() {
  remove_jobs
  /etc/init.d/cron reload
}

restart() {
  start
}

reload() {
  restart
}

status() {
  grep -q lego_wrapper /etc/crontabs/root && echo "running" || echo "inactive" && false
}
