# packages/dnsmasq-dhcp-boot/data/etc/init.d/dnsmasq-dhcp-boot
#!/bin/sh /etc/rc.common

START=14
EXTRA_COMMANDS="status"
EXTRA_HELP="        status          Service status"

config_get() {
  uci get "$1" 2>/dev/null || echo "$2"
}

get_dnsmasqconfdir() {
  config_get "dhcp.@dnsmasq[0].confdir" "/tmp/dnsmasq.d"
}

start() {
  dnsmasqconfdir="$(get_dnsmasqconfdir)"
  [ ! -d "$dnsmasqconfdir" ] && mkdir -p "$dnsmasqconfdir"
  cp -f /etc/dnsmasq_dhcp_boot.conf "$dnsmasqconfdir/dnsmasq_dhcp_boot.conf"
  if [ "$(/etc/init.d/dnsmasq status)" = "running" ]; then
    echo "restart dnsmasq"
    /etc/init.d/dnsmasq restart
  fi
}

stop() {
  dnsmasqconfdir="$(get_dnsmasqconfdir)"
  [ ! -d "$dnsmasqconfdir" ] && mkdir -p "$dnsmasqconfdir"
  rm -f "$dnsmasqconfdir/dnsmasq_dhcp_boot.conf"
  if [ "$(/etc/init.d/dnsmasq status)" = "running" ]; then
    echo "restart dnsmasq"
    /etc/init.d/dnsmasq restart
  fi
}

restart() {
  start
}

reload() {
  restart
}

status() {
  if [ "$(/etc/init.d/dnsmasq status)" = "running" ]; then
    dnsmasqconfdir="$(get_dnsmasqconfdir)"
    [ -f "$dnsmasqconfdir/dnsmasq_dhcp_boot.conf" ] && echo "running" || echo "inactive" && false
  else
    echo "inactive" && false
  fi
}
